# payload_client/systemd/kiosk-storage-setup.service
[Unit]
Description=Kiosk Persistent USB Storage Setup
DefaultDependencies=no
After=local-fs.target
Before=offline-backend.service

[Service]
Type=oneshot
ExecStart=/opt/offline-client/kiosk-storage-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
